# deploy/bread-paper.service
[Unit]
Description=Bread Trading Bot — paper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bread
Group=bread
WorkingDirectory=/home/bread/bread
ExecStart=/home/bread/bread/.venv/bin/python -m bread run --mode paper --dashboard-port 8050
Restart=on-failure
RestartSec=30
TimeoutStopSec=60

# Secrets
EnvironmentFile=/home/bread/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bread-paper

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/bread
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
